# deploy/githubaltmanager.service
# /etc/systemd/system/githubaltmanager.service
# GitHub Alt Manager 后端 Go 服务
#
# 启用:  sudo systemctl daemon-reload
#        sudo systemctl enable --now githubaltmanager
# 查看日志: journalctl -u githubaltmanager -f

[Unit]
Description=GitHub Alt Manager Backend (Go)
After=network.target

[Service]
Type=simple
User=gam
Group=gam
WorkingDirectory=/opt/githubaltmanager/backend
EnvironmentFile=/opt/githubaltmanager/backend/.env
ExecStart=/opt/githubaltmanager/backend/bin/server
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# 安全加固
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/githubaltmanager/backend/data

[Install]
WantedBy=multi-user.target
